# conf/muzzleyconn_openwrt.init
#!/bin/sh /etc/rc.common

START=99

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

EXTRA_COMMANDS="status"
EXTRA_HELP="        status  Returns the status of the daemon"

start() {
        cd /etc/muzzleyconn
        service_start /usr/bin/muzzleyconn
}

stop() {
        service_stop /usr/bin/muzzleyconn
        rm -f /tmp/run/muzzleyconn.pid
}

status() {
        if [ -f /tmp/run/muzzleyconn.pid ];
        then
                echo "Running"
        else
                echo "Stopped"
        fi
}

